# sleep.service
[Unit]
Description=demonstrate simple type

[Service]
Type=simple
ExecStart=sleep 3600
