# stress2.service
[Unit]
Description=Discard unused blocks

[Service]
Type=simple
ExecStart=/usr/bin/dd if=/dev/zero of=/dev/null
CPUShares=2048
